# .config/supervisor/uwsgi.conf
[program:uwsgi]
command=/root/.pyenv/versions/prac_teamproject/bin/uwsgi --ini /srv/00_practice/.config/uwsgi/uwsgi-app.ini

